# QSDM/deploy/systemd/qsdm.service.d/task-registry.conf
[Service]
Environment="QSDM_TASK_REGISTRY_PATH=/opt/qsdm/qsdm-hive-tasks.json"
Environment="QSDM_TASK_ACTION_LOG_PATH=/opt/qsdm/qsdm_task_actions.ndjson"
